# biglinux-improve-compatibility/usr/lib/systemd/system/biglinux-improve-compatibility.service
[Unit]
Description=Using vmtouch to preload programs

# Ordering
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/biglinux-improve-compatibility
User=root
Group=root

[Install]
WantedBy=multi-user.target
